# agents/linux/acc-agent.service
[Unit]
Description=ACC Monitor Agent
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/acc-agent
ExecStart=/usr/bin/python3 /opt/acc-agent/acc_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
